# brewctrl.service
[Unit]
Description=BrewCtrl

[Service]
ExecStart=/opt/brewctrl/brewctrl
WorkingDirectory=/opt/brewctrl


[Install]
WantedBy=multi-user.target
