# brewctrl.service
[Unit]
Description=BrewCtrl

[Service]
ExecStart=/opt/BrewCtrl/brewctrl
WorkingDirectory=/opt/BrewCtrl


[Install]
WantedBy=multi-user.target
